# hostapd-src/wpa-2.6/.pc/wpa_supplicant_fix-dependency-odering-when-invoked-with-dbus.patch/wpa_supplicant/systemd/wpa_supplicant.service.in
[Unit]
Description=WPA supplicant
Before=network.target
Wants=network.target

[Service]
Type=dbus
BusName=@DBUS_INTERFACE@
ExecStart=@BINDIR@/wpa_supplicant -u -s -O /run/wpa_supplicant

[Install]
WantedBy=multi-user.target
Alias=dbus-@DBUS_INTERFACE@.service
